# scripts/production/etc/systemd/system/docker.service.d/production.conf
[Unit]
After=var-lib-docker.mount
BindsTo=var-lib-docker.mount
