# kubernetes/coreos/libvirt/systemd/libvirtd.service.d/99-wsl.conf
[Unit]

After=var-lib-libvirt-images.mount
Requires=var-lib-libvirt-images.mount
