# app/laravel-demo/patch/.docker-rootless/s6/scheduler/finish
#!/command/execlineb -S0

/command/s6-svscanctl -t /app/laravel-docker/.docker-rootless/s6
