# config/s6/services.d/.s6-svscan/finish
#!/bin/sh

set -x

rm -rf $(find /etc/services.d -name 'event')
rm -rf $(find /etc/services.d -name 'supervise')
rm -rf $(find /etc/services.d -name 'lock')
