# kubernetes/systemd/system/mount/tmp-t.mount
[Mount]
# http://www.jinbuguo.com/systemd/systemd.mount.html#

Where=/tmp/t
What=/tmp
Type=none
Options=bind
