# config/s6/s6-overlay/s6-rc.d/php-fpm/finish
#!/command/execlineb -S0

/run/s6/basedir/bin/halt
