# config/s6/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P

#!/usr/bin/env sh

# if test -f ./setup; then
#     source ./setup
# fi

with-contenv

php-fpm -F -R
