# app/laravel-demo/patch/.docker-rootless/s6/.s6-svscan/finish
#!/bin/sh

rm -rf $(find /app/laravel-docker/.docker-rootless/s6 -name 'event')
