# config/s6/laravel/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

exec php /app/laravel/artisan queue:work redis --sleep=3 --tries=3
