# config/s6/services.d/php-fpm/finish
#!/command/execlineb -S0

/command/s6-svscanctl -t /var/run/s6/services
