# config/s6/services.d/php-fpm/run
#!/usr/bin/env sh

if test -f ./setup; then
    source ./setup
fi

exec php-fpm -F
