# config/s6/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

if test -f ./setup; then
    source ./setup
fi

exec php-fpm -F
