# app/laravel-demo/.patch/.docker/s6/.s6-svscan/finish
#!/bin/sh

# rm -rf $(find /etc/services.d -name 'event')
