# app/laravel-demo/.patch/.docker/s6/scheduler/run
#!/usr/bin/env sh

cp /etc/services.d/scheduler/crontab /var/spool/cron/crontabs/root

exec /usr/sbin/crond -fS
