# config/s6/s6-overlay/s6-rc.d/laravel-horizon/run
#!/command/execlineb -P

#!/usr/bin/env sh

# if test -f ./setup; then
#     source ./setup
# fi

php /app/laravel/artisan horizon
