# app/laravel-demo/.patch/.docker/s6/php-fpm/run
#!/usr/bin/env sh

exec /usr/local/sbin/php-fpm -F -R
