# config/s6/demo/s6-overlay/s6-rc.d/myapp/run
# Dropping privileges
#!/usr/bin/execlineb -P
s6-setuidgid daemon
myservice

# Dropping privileges
#!/bin/sh
exec s6-setuidgid daemon myservice

#!/command/execlineb -P
nginx -g "daemon off;"
