# app/laravel-demo/laravel-patch/.docker-rootless/s6/php-fpm/finish
#!/usr/bin/execlineb -S0

s6-svscanctl -t /app/laravel-docker/.docker-rootless/s6
