# config/s6/s6-overlay/s6-rc.d/laravel-queue/run
#!/command/execlineb -P

#!/usr/bin/env sh

# if test -f ./setup; then
#     source ./setup
# fi

with-contenv

php /app/laravel/artisan queue:work redis --sleep=3 --tries=3
