# kubernetes/systemd/system/kube-.service.d/99-default.conf
[Unit]
# ConditionVirtualization=wsl
# ConditionPathIsMountPoint=/wsl/wsl-k8s-data
# RequiresMountsFor=/wsl/wsl-k8s-data
