# kubernetes/systemd/system/kube-.service.d/99-refuse-manual.conf
[Unit]
RefuseManualStart=yes
RefuseManualStop=yes
