# config/s6/laravel-horizon/run
#!/bin/sh

if test -f ./setup; then
    source ./setup
fi

exec php /app/laravel/artisan horizon
