# app/laravel-demo/.patch/.docker-rootless/s6/php-fpm/run
#!/usr/bin/env sh

/usr/local/sbin/php-fpm -t

exec /usr/local/sbin/php-fpm -F
