# kubernetes/systemd/system/path/tmp-t.path
[Path]
# http://www.jinbuguo.com/systemd/systemd.path.html#

PathExists=/tmp/t
Unit=t.service
