# config/s6/services.d/laravel-queue/run
#!/usr/bin/env sh

if test -f ./setup; then
    source ./setup
fi

mkdir -p /var/log/laravel
chown nobody:nogroup /var/log/laravel

exec php /app/laravel/artisan queue:work redis --sleep=3 --tries=3
