# kubernetes/systemd/system/cri-dockerd.service
[Service]
Description=""
Documentation="https://github.com/Mirantis/cri-dockerd"
