# config/s6/s6-overlay/s6-rc.d/crond-log/run
#!/bin/sh
exec logutil-service /var/log/crond
