# kubernetes/systemd/system/systemd-journald.service.d/99-mount.conf
[Unit]
Requires=var-log-journal.mount
After=var-log-journal.mount
Wants=var-log-journal.mount
