# wsl2/supervisor.d/docker.ini
[program:docker]
command=dockerd --data-root=/wsl/wsl-k8s-data/k8s/var/lib/docker --host tcp://0.0.0.0:2376
stdout_logfile=/c/Users/90621/.khs1994-docker-lnmp/wsl-k8s/log/dockerd-stdout.log
stderr_logfile=/c/Users/90621/.khs1994-docker-lnmp/wsl-k8s/log/dockerd-stdout.log
directory=/
autostart=false
autorestart=false
startretries=2
user=root
startsecs=10
