# app/laravel-demo/laravel-patch/.docker/s6/queue/finish
#!/usr/bin/execlineb -S0

s6-svscanctl -t /etc/services.d
