# app/laravel-demo/laravel-patch/.docker/s6/queue/run
#!/usr/bin/env sh

cd /app/laravel-docker && exec php artisan queue:work
